# crontab/task
*/2 * * * *	cd ~/o2gym/crontab && sh refreshAccessToken.sh
